# supervisor/supervisord.conf
[supervisord]
nodaemon=true

[program:updater]
command=/root/apps/updaterules.py start

[program:snort]
command=/root/apps/run.sh
